# device_cio_free.conf
#
# free all devices on startup
#

start on starting rcS

task

console output

exec /sbin/device_cio_free
